# claudeconf/blog.service
[Unit]
Description=Blog Application Service
Documentation=https://github.com/your-repo/blog
After=docker.service mysql.service
Requires=docker.service
Wants=mysql.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/blog/source
Environment=DOCKER_COMPOSE_FILE=docker-compose.prod.yml
Environment=ENV_FILE=.env.prod

# 启动命令
ExecStart=/usr/bin/docker compose --env-file ${ENV_FILE} -f ${DOCKER_COMPOSE_FILE} up -d

# 停止命令
ExecStop=/usr/bin/docker compose --env-file ${ENV_FILE} -f ${DOCKER_COMPOSE_FILE} down

# 重载命令
ExecReload=/usr/bin/docker compose --env-file ${ENV_FILE} -f ${DOCKER_COMPOSE_FILE} restart

# 服务配置
User=root
Group=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=blog-app

# 重启策略
Restart=on-failure
RestartSec=30

# 安全配置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/blog

# 资源限制
MemoryLimit=2G
TasksMax=1000

[Install]
WantedBy=multi-user.target