# claudeconf/blog-backup.service
[Unit]
Description=Blog Database and Files Backup Service
Documentation=https://github.com/your-repo/blog
After=mysql.service blog.service
Requires=mysql.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/blog/scripts

# 备份脚本
ExecStart=/opt/blog/scripts/backup.sh

# 环境变量
Environment=BACKUP_DIR=/opt/blog/data/backups
Environment=MYSQL_USER=root
Environment=DB_NAME=blog_db
Environment=UPLOAD_DIR=/opt/blog/data/uploads
Environment=RETENTION_DAYS=30

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=blog-backup

# 安全配置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/blog/data/backups /opt/blog/logs
ReadOnlyPaths=/opt/blog/data/uploads

# 资源限制
MemoryLimit=512M
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target